# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==================================================================
# Init script for Obsidian Home Assistant add-on
# ==================================================================

# Declare variables
PUID=$(bashio::config 'PUID')
PGID=$(bashio::config 'PGID')
TZ=$(bashio::config 'TZ')

bashio::log.info "Starting Obsidian add-on initialization"
bashio::log.info "Using PUID=$PUID, PGID=$PGID, TZ=$TZ"

# Create addon_config folder
if [ ! -d /addon_config/obsidian ]; then
    bashio::log.info "Creating /addon_config/obsidian"
    mkdir -p /addon_config/obsidian
    chown -R "$PUID:$PGID" /addon_config/obsidian
fi
rm -rf /config
ln -s /addon_config/obsidian /config
bashio::log.info "Config folder linked"

# Create data folder for vault
if [ ! -d /data/obsidian ]; then
    bashio::log.info "Creating /data/obsidian"
    mkdir -p /data/obsidian
    chown -R "$PUID:$PGID" /data/obsidian
fi
rm -rf /vault
ln -s /data/obsidian /vault
bashio::log.info "Vault folder linked"

# Set timezone
export TZ="$TZ"
bashio::log.info "Timezone set to $TZ"

# ==================================================================
# Start Obsidian service
# ==================================================================
exec /usr/bin/obsidian
